# legacy/launchd/com.tyler.speech-transcriber.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tyler.speech-transcriber</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/tyler/p2/speech-transcriber/start.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tyler/p2/speech-transcriber</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/tyler/p2/speech-transcriber/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tyler/p2/speech-transcriber/launchd.error</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>